# components/example/component/applications/prudence-example/crontab
* * * * * /tasks/scheduled/
* * * * * <% application.getSubLogger('scheduled').info('This task will also repeat itself every minute') %>
